# OSMC_GPIO/etc/systemd/system/gpioset.service
[Unit]
Description=Initial setup GPIO

[Service]
Type=idle
ExecStart=/usr/bin/python /home/osmc/gpio.py set

[Install]
WantedBy=multi-user.target
